# shairport_sync/rootfs/etc/s6-overlay/s6-rc.d/svc-shairport-sync/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Addon: Shairport Sync
# Starts the Shairport Sync service
# ==============================================================================

declare CONFIG_FILE=/tmp/shairport-sync.conf
declare USER_CONFIG=/config/shairport-sync.conf

# ---------------------------------------------------------------------------
# Config file: use the user's file if present, otherwise generate from options
# ---------------------------------------------------------------------------
if [[ -f "${USER_CONFIG}" ]]; then
    bashio::log.info "Using custom config from ${USER_CONFIG}"
    cp "${USER_CONFIG}" "${CONFIG_FILE}"
else
    bashio::log.info "Generating config from add-on options..."

    {
        printf 'general = {\n'
        printf '  name = "%s";\n' "$(bashio::config 'airplay_name')"
        printf '  interpolation = "%s";\n' "$(bashio::config 'interpolation')"
        printf '  volume_control_profile = "%s";\n' "$(bashio::config 'volume_control_profile')"
        printf '  log_verbosity = %s;\n' "$(bashio::config 'log_verbosity')"
        printf '  audio_decoded_buffer_desired_length_in_seconds = %s;\n' \
            "$(bashio::config 'audio_decoded_buffer_desired_length_in_seconds')"
        printf '  resync_threshold_in_seconds = %s;\n' \
            "$(bashio::config 'resync_threshold_in_seconds')"
        printf '};\n\n'

        printf 'sessioncontrol = {\n'
        printf '  allow_session_interruption = "%s";\n' "$(bashio::config 'allow_session_interruption')"
        printf '  session_timeout = 120;\n'
        printf '};\n\n'

        printf 'alsa = {\n'
        printf '  audio_backend_buffer_desired_length_in_seconds = %s;\n' \
            "$(bashio::config 'audio_backend_buffer_desired_length_in_seconds')"
        printf '};\n\n'

        printf 'metadata = {\n'
        printf '  enabled = "yes";\n'
        printf '  include_cover_art = "yes";\n'
        printf '};\n'

        if bashio::config.true 'mqtt_enabled'; then
            printf '\nmqtt = {\n'
            printf '  enabled = "yes";\n'
            printf '  hostname = "%s";\n' "$(bashio::config 'mqtt_hostname')"
            printf '  port = %s;\n' "$(bashio::config 'mqtt_port')"
            printf '  topic = "%s";\n' "$(bashio::config 'mqtt_topic')"
            printf '  publish_parsed_metadata = "yes";\n'
            printf '  publish_cover_art = "yes";\n'
            if bashio::config.has_value 'mqtt_username'; then
                printf '  username = "%s";\n' "$(bashio::config 'mqtt_username')"
            fi
            if bashio::config.has_value 'mqtt_password'; then
                printf '  password = "%s";\n' "$(bashio::config 'mqtt_password')"
            fi
            printf '};\n'
        fi
    } > "${CONFIG_FILE}"
fi

# ---------------------------------------------------------------------------
# Wait for Avahi to be fully ready (socket appears after daemon initialises)
# ---------------------------------------------------------------------------
bashio::log.info "Waiting for Avahi to be ready..."
until [[ -S /var/run/avahi-daemon/socket ]]; do
    sleep 0.5
done

# ---------------------------------------------------------------------------
# Wait for nqptp to be ready (creates /dev/shm/nqptp shared memory segment)
# ---------------------------------------------------------------------------
bashio::log.info "Waiting for nqptp to be ready..."
declare -i nqptp_wait=0
until [[ -f /dev/shm/nqptp ]]; do
    if [[ ${nqptp_wait} -ge 30 ]]; then
        bashio::log.warning "nqptp did not become ready after 30 seconds — AirPlay 2 may not work"
        break
    fi
    sleep 1
    nqptp_wait+=1
done
if [[ -f /dev/shm/nqptp ]]; then
    bashio::log.info "nqptp is ready"
fi

# ---------------------------------------------------------------------------
# Start Shairport Sync
# ---------------------------------------------------------------------------
bashio::log.info "Starting Shairport Sync..."
exec shairport-sync -c "${CONFIG_FILE}"
