# shairport_sync/rootfs/etc/s6-overlay/s6-rc.d/svc-shairport-sync/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Addon: Shairport Sync
# Starts the Shairport Sync service
# ==============================================================================

declare airplay_name

# Wait for Avahi to become available
bashio::log.info "Waiting for Avahi to be ready..."
sleep 3

# Determine AirPlay device name
if bashio::config.has_value 'airplay_name'; then
    airplay_name=$(bashio::config 'airplay_name')
else
    airplay_name=$(bashio::info.hostname)
fi

bashio::log.info "Starting Shairport Sync as '${airplay_name}'..."
exec shairport-sync -o pa -m avahi -a "${airplay_name}" -v
