# shairport_sync/rootfs/etc/s6-overlay/s6-rc.d/svc-nqptp/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Addon: Shairport Sync
# Starts the nqptp daemon (required for AirPlay 2)
# ==============================================================================

bashio::log.info "Starting nqptp..."
exec nqptp
