# rootfs/etc/s6-overlay/s6-rc.d/init-dbus/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Addon: Shairport Sync
# Ensures we have a unique D-Bus ID
# ==============================================================================

bashio::log.info "Generating D-Bus UUID..."
dbus-uuidgen --ensure || bashio::exit.nok 'Failed to generate a unique D-Bus ID'
mkdir -p /var/run/dbus
