# rootfs/etc/s6-overlay/s6-rc.d/svc-avahi/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Home Assistant Addon: Shairport Sync
# Start the Avahi daemon
# ==============================================================================

# Wait until D-Bus is available
bashio::log.info "Waiting for D-Bus..."
until [[ -e /var/run/dbus/system_bus_socket ]]; do
  sleep 0.5
done

# Start the Avahi daemon
bashio::log.info "Starting Avahi daemon..."
if bashio::debug; then
  exec avahi-daemon --debug
else
  exec avahi-daemon
fi
